# srcs/InfluxDB/supervisord.conf
[supervisord]
nodaemon=true

[program:influxd]
command=influxd
autorestart=true

[program:telegraf]
command=telegraf --config telegraf.conf
autorestart=true