# srcs/Mysql/supervisord.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=mysqld -u root --datadir='/var/lib/mysql'
autorestart=true

[program:telegraf]
command=telegraf --config telegraf.conf
autorestart=true