# srcs/Nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx
autorestart=true

[program:sshd]
command=/usr/sbin/sshd

[program:telegraf]
command=telegraf --config telegraf.conf
autorestart=true