# srcs/ftps/supervisord.conf
[supervisord]
nodaemon=true

[program:vsftpd]
command=vsftpd /etc/vsftpd.conf
autorestart=true

[program:telegraf]
command=telegraf --config telegraf.conf
autorestart=true